# deploy/etc/cron.d/lims2_global
# LIMS2 全局(不分 lab)的 crontab
# 修改该文件时, 应 find 各 labs 中是否有重载, 并对重载也做修改
0 * * * * root /etc/lims2/backup 2>&1
